# launchd/com.user.ai-trading-agent.scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.ai-trading-agent.scan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/nishazdhuka/Robinhood AI/scripts/loop-tick.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/nishazdhuka/Robinhood AI</string>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>/Users/nishazdhuka/Robinhood AI/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nishazdhuka/Robinhood AI/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>

