# launchd/com.user.ai-trading-agent.scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.ai-trading-agent.scan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/nishazdhuka/Robinhood AI/scripts/run_daily_scan.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/nishazdhuka/Robinhood AI</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>45</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>45</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>45</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>45</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>45</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/nishazdhuka/Robinhood AI/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nishazdhuka/Robinhood AI/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>

